# config/supervisor.conf
[supervisord]
nodaemon=true

[program:php-fpm]
command=/usr/sbin/php5-fpm -F
autostart=true
autorestart=true
priority=5
stdout_events_enabled=true
stderr_events_enabled=true

[include]
files = /etc/supervisor/conf.d/*.conf
